# installer/systemd/serino.service
[Unit]
Description=Serino Application Stack
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/serino
ExecStart=/bin/bash -lc 'set -a; source /etc/serino/serino.env; set +a; export COMPOSE_PROJECT_NAME=serino; exec docker compose -f /opt/serino/docker-compose.release.yml up -d --remove-orphans'
ExecStop=/bin/bash -lc 'set -a; source /etc/serino/serino.env; set +a; export COMPOSE_PROJECT_NAME=serino; exec docker compose -f /opt/serino/docker-compose.release.yml down'
ExecReload=/bin/bash -lc 'set -a; source /etc/serino/serino.env; set +a; export COMPOSE_PROJECT_NAME=serino; exec docker compose -f /opt/serino/docker-compose.release.yml up -d --remove-orphans'
TimeoutStartSec=0
TimeoutStopSec=180

[Install]
WantedBy=multi-user.target
